# templates/default/init/stud.conf.erb
# stub
#
# Stub provides high speed SSL termination

description "Stub Server"

start on runlevel [2345]
stop on runlevel [!2345]

#respawn
#respawn limit 10 5

#umask 007

script
for configFile in /etc/stud/configs-enabled/*.conf
do
   exec stud --config=$configFile
done
end script

console output