# templates/default/init/stud.conf.erb
# stub
#
# Stub provides high speed SSL termination

description "Stub Server"

start on runlevel [2345]
stop on runlevel [!2345]

#respawn
#respawn limit 10 5

#umask 007

script
for configFile in /etc/stub/configs-enabled/*.conf
do
   exec stub --config=$configFile
done
end script

exec stub
console output